# supervisord.conf
[inet_http_server]
port = :9001

[program:rig]
directory = /home/download/
command = http-server .